# src/Nemu/Installer/linux/debian/postinst.in
#!/bin/sh
#
# Copyright (C) 2006-2015 Oracle Corporation
#
# This file is part of VirtualBox Open Source Edition (OSE), as
# available from http://www.virtualbox.org. This file is free software;
# you can redistribute it and/or modify it under the terms of the GNU
# General Public License (GPL) as published by the Free Software
# Foundation, in version 2 as it comes in the "COPYING" file of the
# VirtualBox OSE distribution. VirtualBox OSE is distributed in the
# hope that it will be useful, but WITHOUT ANY WARRANTY of any kind.
#

# we can be called with the following arguments (6.5 of Debian policy):
#  configure:         (our version): installing/configuring new version
#  abort-upgrade:     (old version): upgrading to a new version failed
#  abort-remove:      (our version): removing this package failed
#  abort-deconfigure: (our version): error during resolving conflicts

#include installer-common.sh

LOG="/var/log/nemu-install.log"

# defaults
[ -r /etc/default/virtualbox ] && . /etc/default/virtualbox

if [ "$1" = "configure" ]; then

  # for debconf
  . /usr/share/debconf/confmodule
  db_version 2.0

  # remove old cruft
  if [ -f /etc/init.d/nemudrv.sh ]; then
    echo "Found old version of /etc/init.d/nemudrv.sh, removing."
    rm /etc/init.d/nemudrv.sh
    update-rc.d nemudrv.sh remove >/dev/null
  fi
  if [ -f /etc/nemu/nemu.cfg ]; then
    echo "Found old version of /etc/nemu/nemu.cfg, removing."
    rm /etc/nemu/nemu.cfg
  fi

  # create users groups (disable with INSTALL_NO_GROUP=1 in /etc/default/virtualbox)
  if [ "$INSTALL_NO_GROUP" != "1" ]; then
    db_input low virtualbox/group-nemuusers || true
    db_go || true
    addgroup --system nemuusers || true
  fi

  # install udev rule (disable with INSTALL_NO_UDEV=1 in
  # /etc/default/virtualbox) and /dev/nemudrv and /dev/nemuusb/*/* device nodes
  install_device_node_setup root 0600 /usr/share/virtualbox "${usb_group}"

  # The starters need to be Suid root. They drop the privileges before starting
  # the real frontend.
  if ! dpkg-statoverride --list /usr/lib/virtualbox/VirtualBox > /dev/null 2>&1; then
    chmod 4511 /usr/lib/virtualbox/VirtualBox
  fi
  if ! dpkg-statoverride --list /usr/lib/virtualbox/NemuHeadless > /dev/null 2>&1; then
    chmod 4511 /usr/lib/virtualbox/NemuHeadless
  fi
  if ! dpkg-statoverride --list /usr/lib/virtualbox/NemuSDL > /dev/null 2>&1; then
    chmod 4511 /usr/lib/virtualbox/NemuSDL
  fi
  if ! dpkg-statoverride --list /usr/lib/virtualbox/NemuNetDHCP > /dev/null 2>&1; then
    chmod 4511 /usr/lib/virtualbox/NemuNetDHCP
  fi
  if ! dpkg-statoverride --list /usr/lib/virtualbox/NemuNetNAT > /dev/null 2>&1; then
    chmod 4511 /usr/lib/virtualbox/NemuNetNAT
  fi
  if ! dpkg-statoverride --list /usr/lib/virtualbox/NemuNetAdpCtl > /dev/null 2>&1; then
    chmod 4511 /usr/lib/virtualbox/NemuNetAdpCtl
  fi
  if [ -x /usr/lib/virtualbox/NemuVolInfo ]; then
    if ! dpkg-statoverride --list /usr/lib/virtualbox/NemuVolInfo > /dev/null 2>&1; then
      chmod 4511 /usr/lib/virtualbox/NemuVolInfo
    fi
  fi

  # if INSTALL_NO_NEMUDRV is set to 1, remove all shipped modules
  if [ "$INSTALL_NO_NEMUDRV" = "1" ]; then
    rm -f /lib/modules/*/misc/nemudrv.ko
    rm -f /lib/modules/*/misc/nemunetflt.ko
    rm -f /lib/modules/*/misc/nemunetadp.ko
  fi

fi # $1 = "configure"

#DEBHELPER#

BUILD_MODULES=0
REGISTER_MODULES=1
# Disable module compilation with INSTALL_NO_NEMUDRV=1 in /etc/default/virtualbox
if [ ! -f /lib/modules/`uname -r`/misc/nemudrv.ko ]; then
  REGISTER_MODULES=0
  if [ "$1" = "configure" ]; then
    if [ "$INSTALL_NO_NEMUDRV" != "1" ]; then
      db_get virtualbox/module-compilation-allowed
      if [ "$RET" = "false" ]; then
        cat << EOF
Unable to find a precompiled module for the current kernel
though module compilation denied by debconf setting.
EOF
      else
        db_input low virtualbox/module-compilation-allowed || true
        db_go || true
        db_get virtualbox/module-compilation-allowed
        if [ "$RET" = "true" ]; then
          BUILD_MODULES=1
        fi
      fi
    fi
  fi # $1 = "configure"
fi

if [ $BUILD_MODULES -eq 1 ]; then
  /usr/lib/virtualbox/nemudrv.sh setup || true
fi

# Install and start the new service scripts.
PRERM_DKMS=
test "${REGISTER_MODULES}" = 1 && PRERM_DKMS="--dkms %VER%"
POSTINST_START=
test "${INSTALL_NO_NEMUDRV}" = 1 && POSTINST_START=--nostart
/usr/lib/virtualbox/prerm-common.sh ${PRERM_DKMS} || true
/usr/lib/virtualbox/postinst-common.sh ${POSTINST_START} > /dev/null || true

exit 0
